# skills/autoresearch-paper/references/launchd/com.mavis.plan-rescue-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mavis.plan-rescue-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>{{RUNTIME_SCRIPTS}}/plan-rescue-daemon.py</string>
        <string>--once</string>
    </array>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>{{LOG_DIR}}/plan-rescue-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_DIR}}/plan-rescue-daemon.err.log</string>

    <key>WorkingDirectory</key>
    <string>{{HOME}}</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{PATH}}</string>
        <key>PYTHONPATH</key>
        <string>{{RUNTIME_SCRIPTS}}</string>
    </dict>
</dict>
</plist>
